# paper_trader/systemd/argus-paper-trader.service
[Unit]
Description=Argus paper trader (micro_deviation maker strategy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/paper-trader
EnvironmentFile=/home/ubuntu/collector-dhan/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/paper-trader/venv/bin/python -m paper_trader.main
Restart=on-failure
RestartSec=15
TimeoutStopSec=60
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
