# paper_trader/systemd/argus-paper-trader-stop.service
[Unit]
Description=Stop Argus paper trader after market close

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop argus-paper-trader.service
